# dot_config/systemd/user/gpu-screen-recorder-replay.service
[Unit]
Description=GPU Screen Recorder (Instant Replay)

[Service]
ExecStart=/home/insprill/bin/start_gpu_recorder_replay.sh
Restart=always

[Install]
WantedBy=default.target
