# bin/setup/etc/systemd/system/uni-sync.timer
[Unit]
Description=Uni-Sync Timer

[Timer]
OnBootSec=1min
OnUnitActiveSec=1min

[Install]
WantedBy=timers.target
